# deploy/systemd/personal-pasteboard-backup.timer
[Unit]
Description=Daily personal pasteboard encrypted backup

[Timer]
OnCalendar=*-*-* 03:20:00
Persistent=true

[Install]
WantedBy=timers.target
